# package/custom-files/files/etc/init.d/blade3-boot
#!/bin/sh /etc/rc.common

START=15
STOP=20

start() {
    /bin/nodectl poweron --all
}

stop() {
    printf "\tThe node is powered on\n"
}
